# supervisor/celeryscheduleflower.ini
[program:celeryscheduleflow]
command=/usr/local/python3.8/bin/celery flower -A celeryschedule worker --address=0.0.0.0 --port=5555  --loglevel=info
directory=/root/celeryproject
stdout_logfile=/root/celeryproject/celeryschedule/celeryflower.log
autorestart=true
redirect_stderr=true
